# systemd/dictate-agent.service
[Unit]
Description=Dictate Agent - Voice-activated AI assistant
Documentation=https://github.com/jakedevar/dictate_agent
After=network.target graphical-session.target
Wants=pulseaudio.service

[Service]
Type=simple
ExecStart=%h/dictate_agent/target/release/dictate-agent
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=RUST_LOG=dictate_agent=info

# Resource limits
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=default.target
